# target-root/etc/update-motd.d/90-testadura
#!/bin/bash
source ../../usr/local/lib/testadura/common/default-colors.sh

# Default gateway (overall)
GW=$(ip route show default 2>/dev/null | awk '/^default/ {print $3; exit}')
DNS=$(resolvectl status 2>/dev/null \
      | awk -F': ' '/DNS Servers/ {print $2}' \
      | tr ' ' '\n' | sort -u | xargs)
[ -n "$DNS" ] || DNS=$(awk '/^nameserver/ {print $2}' /etc/resolv.conf | xargs)

UPTIME=$(uptime -p | sed 's/up //')
MEM_USED=$(free -h | awk '/^Mem:/ {print $3 "/" $2}')
DISK_USED=$(df -h / | awk 'NR==2 {print $3 "/" $2 " (" $5 ")"}')
LOADAVG=$(awk '{print $1", "$2", "$3}' /proc/loadavg)

STAMP="/var/lib/apt/periodic/update-success-stamp"
if [[ -f $STAMP ]]; then
    last_update=$(date -r "$STAMP" +%s)
    now=$(date +%s)
    days=$(( (now - last_update) / 86400 ))
    formatted_date=$(date +"%a, %d %b %Y %H:%M" -r "$STAMP")
    LAST_UPDATE="$formatted_date ($days day$([[ $days -ne 1 ]] && echo s) ago)"
else
    LAST_UPDATE="Never"
    ((days=9999))
fi

# Theme
ttl=$BOLD_YELLOW
lbl=$BOLD_YELLOW
txt=$BOLD_SILVER
dsc=$WHITE
rst=$RESET

printf "======================================================================================\n"
printf " ${ttl}Cloned VM — Ubuntu Server 24.04 LTS$rst\n\n"
printf "    ${lbl}Hostname    :${rst} %s\n" "$(hostname -f)" # Hostname
printf "    ${lbl}IPv4        :${rst} %s\n" "$(hostname -I | awk '{print $1}')" # IP addresses (IPv4 only)
printf "    ${lbl}Gateway     :${rst} %s\n" "${GW:-N/A}"
printf "    ${lbl}DNS         :${rst} %s\n" "${DNS:-N/A}"
printf "    ${lbl}Last update :${rst} %s\n" "${LAST_UPDATE}"

if (( days > 7 )); then
       printf "\n--------------------------------------------------------------------------------------\n"
       if (( days >= 9999 )); then
           printf "\n %sSystem was never updated.%s" "$BOLD_RED" "$RESET"
       else
           printf "\n %sLast update was $days ago.%s" "$ORANGE" "$RESET"
       fi
       printf "\n %sConsider updating by running on of the following commands:${rst}\n" "$BOLD_YELLOW"
       printf "   • sudo testadura.systemupdate.sh full\n"
       printf "     Full upgrade including new dependencies\n"
       printf "   • sudo testadura.systemupdate.sh\n"
       printf "     Package upgrade only\n"
fi
printf "\n--------------------------------------------------------------------------------------\n"
printf " ${ttl}Testadura Consultancy — Post-clone Initialization Scripts${rst}\n\n"
printf " ${lbl}Runs Automatically on First Boot${rst}\n"
printf " ${txt}  • testadura.firstboot.sh${rst}\n"
printf "     ${dsc}Sets up machine ID and prepares the system${rst}\n\n"
printf " ${lbl}Run Manually${rst}\n"
printf " ${txt}  • testadura.netconfig.sh${rst}\n"
printf "     ${dsc}Interactive network configuration${rst}\n"
printf " ${txt}  • testadura.enablessh.sh${rst}\n"
printf "     ${dsc}Regenerates SSH host keys and enables SSH access${rst}\n\n"
printf " ${lbl}Other Components${rst}\n"
printf " ${txt}  • testadura.firstboot.service${rst}\n"
printf "     ${dsc}Systemd service triggering firstboot script${rst}\n"
printf "  ${txt}   • testadura.preparetemplate.sh${rst}\n"
printf "     ${dsc}Prepares the system for cloning as a reusable template${rst}\n"
printf "\n--------------------------------------------------------------------------------------\n"
printf " ${lbl}%-22s %-15s %-18s %-17s${rst}\n" \
       "Uptime" "Memory" "Disk Usage" "CPU Load (1/5/15m)"
printf " %-22.22s %-15.15s %-18.18s %-17.17s\n" \
       "$UPTIME" "$MEM_USED" "$DISK_USED" "$LOADAVG"
printf "\n--------------------------------------------------------------------------------------\n"
printf "\n"
printf "======================================================================================\n"
